# .config/profile.d/bemenu.sh
#!/bin/sh

export BEMENU_BACKEND="wayland"
export BEMENU_OPTS="-i -H 30 --fn 'JetBrainsMono NF 14px' --tf '#1b2224' --hf '#2bb29a' --tb '#2bb29a' --hb '#1b2224' --nb '#1b2224' --fb '#1b2224' --ab '#1b2224'"
